# tclmd.service
[Unit]
Description=The TSClient LEGACY Centralized Lock Manager

[Service]
Type=simple
ExecStart=/usr/sbin/tclmd
